# eos-update-notifier/eos-update-notifier.timer
[Unit]
Description=Run EOS update notifier on a schedule

[Timer]
OnStartupSec=30sec
OnCalendar=daily
 
[Install]
WantedBy=timers.target
