# eos-update-notifier/eos-update-notifier.timer
[Unit]
Description=Run EOS update notifier periodically

[Timer]
OnStartupSec=30sec
OnCalendar=daily
Persistent=true
 
[Install]
WantedBy=timers.target
