# eos-update-notifier/eos-update-notifier.timer
[Unit]
Description=EOS update notifier runs periodically, and optionally after each boot, or if missed.

[Timer]
OnCalendar=weekly
# OnStartupSec=
Persistent=false
RandomizedDelaySec=10min
 
[Install]
WantedBy=timers.target
