# z_archived/reflector-auto/reflector-auto.timer
[Unit]
Description=Run reflector-auto periodically

[Timer]
OnCalendar=weekly
Persistent=true
 
[Install]
WantedBy=timers.target
